# vm/rootfs/etc/init.d/firewall4ai
#!/sbin/openrc-run

description="Firewall4AI AI proxy"
command="/usr/bin/firewall4ai"
command_args="-config /etc/firewall4ai/config.json"
command_background=true
pidfile="/run/${RC_SVCNAME}.pid"
output_log="/var/log/firewall4ai.log"
error_log="/var/log/firewall4ai.log"

depend() {
    need net
    after dnsmasq
    after local
}

start_pre() {
    checkpath -d -m 0755 /var/lib/firewall4ai
}
